# bin/jobs.conf
[program:job_worker]
command=/service/miduoduo/www/miduoduo job-worker/run
autostart = true
